# profile.d/qt5-theme.sh
#Forces QT to use the theme set by Qt5ct
export QT_QPA_PLATFORMTHEME=qt5ct
