# profile.d/qt5-awesome-mouse-fix.sh
#Fixes Awesome panels not working when QT apps are running
export QT_XCB_NO_XI2_MOUSE=1

