# autostart/Compositor (compton).desktop
[Desktop Entry]
Exec=compton -b
Name=Compositor (compton)
Type=Application
Version=1.0
