# autostart/Global Patchbay Canvas.desktop
[Desktop Entry]
Exec=/usr/bin/catia
Name=Global Patchbay Canvas
Type=Application
Version=1.0
