# autostart/JACK MIDI Clock.desktop
[Desktop Entry]
Exec=jack_midi_clock 
Name=JACK MIDI Clock
Type=Application
Version=1.0
