# autostart/Laptop Battery Indicator.desktop
[Desktop Entry]
Exec=fdpowermon
Name=Laptop Battery Indicator
Type=Application
Version=1.0
