# autostart/MIDI ROUTER 2.desktop
[Desktop Entry]
Exec=qmidiroute ~/SYSTEM/GLOBAL/MIDI-ROUTER/MIDI\\ ROUTER\\ 2.qmr
Hidden=true
Name=MIDI ROUTER 2
Type=Application
Version=1.0
