# autostart/MIDI ROUTER 6.desktop
[Desktop Entry]
Exec=qmidiroute ~/SYSTEM/GLOBAL/MIDI-ROUTER/MIDI\\ ROUTER\\ 6.qmr
Hidden=true
Name=MIDI ROUTER 6
Type=Application
Version=1.0
