# autostart/MIDI ROUTER 5.desktop
[Desktop Entry]
Exec=qmidiroute ~/SYSTEM/GLOBAL/MIDI-ROUTER/MIDI\\ ROUTER\\ 5.qmr
Hidden=true
Name=MIDI ROUTER 5
Type=Application
Version=1.0
