# autostart/Network Manager.desktop
[Desktop Entry]
Exec=frisbee-tray
Name=Network Manager
Type=Application
Version=1.0
