# autostart/Non Session Manager.desktop
[Desktop Entry]
Exec=non-session-manager
Name=Non Session Manager
Type=Application
Version=1.0
