# autostart/Terminal (tilda).desktop
[Desktop Entry]
Exec=tilda
Name=Terminal (tilda)
Type=Application
Version=1.0
