# autostart/zNotes.desktop
[Desktop Entry]
Exec=znotes
Name=zNotes
Type=Application
Version=1.0
